# conf/mailer.conf
[program:mailer]
command= /root/Downloads/node-v8.10.0-linux-armv7l/bin/node pduDecode.js
directory=/root/Code/transfer-sms/mailer

[supervisord]
logfile= ../logs/mailer.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=../logs/mailer.pid ; (supervisord pidfile;default supervisord.pid)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=root                 ; (default is current user, required if root)
childlogdir=../logs/supervisord/            ; ('AUTO' child log dir, default $TEMP)
